# packages/collector/service/etm12-collector.service
[Unit]
Description=ETM12 Homedata Collector
Documentation=https://github.com/stuf/homedata-etm12
After=network.target

[Service]
Type=simple
User=homedata
ExecStart=/home/homedata/.etm12/bin/collector.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
